# package/tailscale-install.service
[Unit]
Description=Ensure that Tailscale is installed on your device
After=network.target

[Service]
Type=oneshot
RemainAfterExit=no
Restart=no
Environment=DEBIAN_FRONTEND=noninteractive
ExecStart=/bin/bash /data/tailscale/manage.sh on-boot

[Install]
WantedBy=multi-user.target